# backend/crontab
# Run pipeline for all configured agents every 2 hours
0 */2 * * * python -m app.pipeline run-all >> /proc/1/fd/1 2>> /proc/1/fd/2
